# rootfs/etc/services.d/jira/run
#!/usr/bin/execlineb -P

cd /opt/jira

with-contenv
s6-setuidgid jira
foreground { s6-rmrf /var/lib/jira/.jira-home.lock }
trap -x
{
  term { /opt/jira/bin/stop-jira.sh }
}
/opt/jira/bin/start-jira.sh -fg
